# docker/s6-rc.d/hitster/run
#!/command/with-contenv sh
/hitster/server